# etc/systemd/system/eureka.service
[Unit]
Description=Eureka Blog System

[Install]
WantedBy=multi-user.target

[Service]
User=root
Type=fork
PIDFile=/var/run/eureka/uwsgi.pid
ExecStart=/usr/libexec/eureka/eureka.sh
ExecStop=/bin/kill -QUIT $MAINPID
Restart=always
